# openstack-trove-conductor.service
[Unit]
Description=OpenStack Trove Conductor Service
After=syslog.target network.target mariadb.service

[Service]
Type=simple
User=trove
ExecStart=/usr/bin/trove-conductor

[Install]
WantedBy=multi-user.target

